# mimo-listener-2.service
[Unit]
Description=Mimo Firebase Print Listener (CV-001 / printpi)
After=network.target

[Service]
Type=simple
User=printpi
WorkingDirectory=/home/printpi
Environment="PYTHONUNBUFFERED=1"
Environment="BW_PRINTER_NAME=Brother_IPP"
Environment="COLOR_PRINTER_NAME=Brother_IPP"
Environment="KIOSK_ID=CV-001"
Environment="GRPC_DEFAULT_SSL_ROOTS_FILE_PATH=/etc/ssl/certs/ca-certificates.crt"
ExecStart=/usr/bin/python3 /home/printpi/firebase_listener.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
